# deploy/app.conf
description "Run Docker App Setup"
start on started docker
task
respawn limit 15 5
chdir /home/ubuntu
exec ./setup.sh 2>&1 >> /var/log/app-setup.log
